# intel-undervolt.service.in
[Unit]
Description=Intel Undervolt Service
After=multi-user.target suspend.target hibernate.target hybrid-sleep.target

[Service]
Type=oneshot
ExecStart=%BINDIR%/intel-undervolt apply

[Install]
WantedBy=multi-user.target suspend.target hibernate.target hybrid-sleep.target
