# intel-undervolt-loop.service.in
[Unit]
Description=Intel Undervolt Loop Service
After=multi-user.target

[Service]
Type=simple
ExecStart=%BINDIR%/intel-undervolt daemon
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=multi-user.target
