# intel-undervolt.service
[Unit]
Description=Intel Undervolt Service
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/intel-undervolt apply
